# batman-adv/files/etc/uci-defaults/98-batman-adv-nss-enable
#!/bin/sh

[ -x /etc/init.d/batman-adv-nss ] || exit 0
/etc/init.d/batman-adv-nss enable

exit 0
